# templates/supervisord.ini
# {{ ansible_managed }}
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
#logfile = {{ sketchy_home }}/supervisord.log
logfile = /var/log/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
pidfile = /tmp/supervisord.pid
nodaemon = False
minfds = 1024
minprocs = 200
umask = 022
identifier = supervisor
directory = /tmp
nocleanup = true
childlogdir = /tmp

[program:redis]
command=redis-server

[program:celeryd]
command={{ sketchy_root }}/env-sketchy/bin/celery worker -A sketchy.celery
directory={{ sketchy_root }}
user={{ sketchy_user }}
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600

[program:gunicorn]
command={{ sketchy_root }}/env-sketchy/bin/gunicorn sketchy:app -b {{ sketchy_host }}:{{ sketchy_port }}
directory={{ sketchy_root }}
user={{ sketchy_user }}
autostart=true
autorestart=true
redirect_stderr=true

