# templates/supervisor.ini
[program:celeryd]
command={{ sketchy_root }}/env-sketchy/bin/celery worker -A sketchy.celery
directory={{ sketchy_root }}
user={{ sketchy_user }}
stdout_logfile=/var/log/supervisor/sketchy-celery.log
redirect_stderr=true
autostart=true
autorestart=true
startsecs=10

[program:gunicorn]
command={{ sketchy_root }}/env-sketchy/bin/gunicorn sketchy:app -b {{ sketchy_host }}:{{ sketchy_port }}
directory={{ sketchy_root }}
user={{ sketchy_user }}
stdout_logfile=/var/log/supervisor/sketchy-gunicorn.log
redirect_stderr=true
autostart=true
autorestart=true
startsecs=10

