# templates/supervisor.ini
[program:celeryd]
command=celery worker -A sketchy.celery
directory={{ sketchy_root }}
user={{ sketchy_user }}

[program:gunicorn]
command=python gunicorn sketchy:app -b {{ sketchy_host }}:{{ sketchy_port }}
directory={{ sketchy_root }}
user={{ sketchy_user }}

